# config/files/etc/systemd/system/rpm-ostreed-automatic.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnBootSec=15m
OnUnitInactiveSec=
OnCalendar=*-*-* 17:30:00 UTC
Persistent=true
